# thinkpad-fan-tweak.service.in
[Unit]
Description=ThinkPad fan tweak

[Service]
Type=simple
Restart=no
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@

[Install]
WantedBy=multi-user.target 
